# deploy/funasr-compat-asr.supervisor.conf
[program:funasr-nano-asr]
directory=/root/autodl-tmp/funasr-nano-service
command=/bin/bash -lc 'export MODELSCOPE_CACHE=/root/autodl-tmp/models; export HF_HOME=/root/autodl-tmp/hf-home; export TRANSFORMERS_CACHE=/root/autodl-tmp/hf-home/transformers; export CUDA_VISIBLE_DEVICES=0; exec /root/autodl-tmp/funasr-nano-venv/bin/python -u serve_realtime_ws_compat.py --host 127.0.0.1 --port 10095 --model paraformer-zh --device cuda:0'
autostart=true
autorestart=true
startsecs=10
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/root/autodl-tmp/logs/funasr-compat-asr.log
stderr_logfile=/root/autodl-tmp/logs/funasr-compat-asr.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=PYTHONUNBUFFERED="1"
