# deploy/phdebate-target-12234.supervisor.conf
[program:jixia-debate]
priority=10
directory=/home/ubuntu/sunsq/debateall/apps/backend
command=/bin/bash -lc 'set -a; source /home/ubuntu/sunsq/debateall/.env; set +a; exec /home/ubuntu/sunsq/debateall/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 12234 --no-access-log'
environment=PYTHONPATH="/home/ubuntu/sunsq/debateall/apps/backend",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/home/ubuntu/sunsq/debateall/runtime/logs/jixia-debate.stdout.log
stderr_logfile=/home/ubuntu/sunsq/debateall/runtime/logs/jixia-debate.stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:jixia-voice-agent]
priority=20
directory=/home/ubuntu/sunsq/debateall/apps/voice_agent
command=/bin/bash -lc 'set -a; source /home/ubuntu/sunsq/debateall/.env; set +a; exec /home/ubuntu/sunsq/debateall/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 12235 --no-access-log'
environment=PYTHONPATH="/home/ubuntu/sunsq/debateall/apps/voice_agent",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=5
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/home/ubuntu/sunsq/debateall/runtime/logs/voice-agent.stdout.log
stderr_logfile=/home/ubuntu/sunsq/debateall/runtime/logs/voice-agent.stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:jixia-qwen3-tts]
priority=40
directory=/home/ubuntu/sunsq/debateall/services/qwen3-tts-openai
command=/home/ubuntu/sunsq/debateall/services/qwen3-tts-openai/start_after_asr.sh
environment=PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=20
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/home/ubuntu/sunsq/debateall/runtime/logs/qwen3-tts.stdout.log
stderr_logfile=/home/ubuntu/sunsq/debateall/runtime/logs/qwen3-tts.stderr.log
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:jixia-funasr-nano-asr]
priority=50
directory=/home/ubuntu/sunsq/debateall/services/funasr-nano-service
command=/bin/bash -lc 'source /home/ubuntu/sunsq/debateall/services/funasr-nano-venv/bin/activate; export MODELSCOPE_CACHE=/home/ubuntu/sunsq/debateall/services/models; export HF_HOME=/home/ubuntu/sunsq/debateall/services/hf-home; export TRANSFORMERS_CACHE=/home/ubuntu/sunsq/debateall/services/hf-home/transformers; export CUDA_VISIBLE_DEVICES=${CUDA_VISIBLE_DEVICES:-0}; exec python -u serve_realtime_ws_compat.py --host 127.0.0.1 --port 10095 --model paraformer-zh --device cuda:0'
environment=PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startsecs=20
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/home/ubuntu/sunsq/debateall/runtime/logs/funasr-nano-asr.stdout.log
stderr_logfile=/home/ubuntu/sunsq/debateall/runtime/logs/funasr-nano-asr.stderr.log
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB
stdout_logfile_backups=5
stderr_logfile_backups=5
